# dot_config/systemd/user/llama-server@.service
[Unit]
Description=llama-server (%i)
After=network.target

[Service]
EnvironmentFile=%h/.config/llama-server/%i.env
ExecStart=%h/.local/src/llama.cpp/build-stable/bin/llama-server \
  -m ${MODEL} -ngl ${NGL} --device ${DEVICE} \
  -fa ${FA} -c ${CTX} -b ${BATCH} -ub ${UBATCH} \
  --host ${HOST} --port ${PORT} \
  --reasoning ${REASONING_MODE} \
  --temp ${TEMP} \
  --top-p ${TOP_P} \
  --top-k ${TOP_K} \
  --min-p ${MIN_P} \
  --repeat-penalty ${REPEAT_PENALTY} \
  --presence-penalty ${PRESENCE_PENALTY}
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=default.target
