# dot_config/systemd/user/llama-update.timer
[Unit]
Description=Weekly llama.cpp update

[Timer]
OnCalendar=Sun 04:00
Persistent=true
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
