# dot_config/systemd/user/llama-update.service
[Unit]
Description=Update and build llama.cpp (build-latest)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/llama-update.sh
